# crontab
0	22	*	*	*	/var/www/application/vendor/bin/laminas update-command -t 24
